# roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User=node_exporter
ExecStart=/usr/bin/node_exporter --web.listen-address={{ wireguard_ip }}:9100

[Install]
WantedBy=multi-user.target
