# roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter

[Service]
User=node_exporter
ExecStart=/usr/bin/node_exporter --web.listen-address=127.0.0.1:9100

[Install]
WantedBy=multi-user.target
